# execution_engine/scheduler/execution-engine.timer.example
[Unit]
Description=Run Crypto Engine execution cycle every 5 minutes

[Timer]
OnCalendar=*:01/5
AccuracySec=1s
RandomizedDelaySec=0
Unit=execution-engine.service

[Install]
WantedBy=timers.target
